# root/etc/services.d/cleaning/finish
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

pids="$(ps -ef | grep 'cleanup.sh' | head -n 1 | grep -v grep | awk '{print $1}')"
if [ "$pids" != "" ]; then
  kill -15 $pids
  sleep 2
fi
